# config/runit/sv/agetty-tty1/run
#!/bin/sh
exec agetty -a ert -L tty1 linux
