# config/runit/sv/agetty-tty2/run
#!/bin/sh
exec setsid agetty -L tty2 linux
