# config/runit/sv/wpa_supplicant/run
#!/bin/sh
IFACE=$(ip link | grep -E "wl[a-z0-9]+" -o | head -1)
[ -z "$IFACE" ] && { sleep 5; exit 1; }
exec wpa_supplicant -i "$IFACE" -c /etc/wpa_supplicant/wpa_supplicant.conf -s
