# config/etc/profile.d/locale.sh
# UTF-8 locale (required for tmux and other modern tools)
export LANG=en_US.UTF-8
