# config/etc/profile.d/pk.sh
# pk package manager settings
# PK_ROOT="/" - system install (absolute paths)
# PK_ROOT="/mnt/lfs/" - chroot install
# PK_ROOT="" - local/safe mode (relative paths)
export PK_ROOT=/
