# config/runit/sv/dhcpcd/run
#!/bin/sh
# runit service: dhcpcd DHCP client

exec dhcpcd -B -M
