# config/runit/sv/sshd/run
#!/bin/sh
# runit service: OpenSSH daemon

# Generate host keys if missing
if [ ! -f /etc/ssh/ssh_host_rsa_key ]; then
    ssh-keygen -A
fi

exec /usr/sbin/sshd -D -e
