# config/runit/sv/iwd/run
#!/bin/sh
# runit service: iwd wireless daemon

sv check dbus >/dev/null || exit 1
exec iwd
