# config/runit/sv/agetty-tty5/run
#!/bin/sh
exec setsid agetty -L tty5 linux
