# config/runit/sv/agetty-tty3/run
#!/bin/sh
exec setsid agetty -L tty3 linux
