# config/runit/sv/agetty-tty4/run
#!/bin/sh
exec setsid agetty -L tty4 linux
