# config/runit/sv/agetty-tty6/run
#!/bin/sh
exec setsid agetty -L tty6 linux
